# 各版本的files/人像探针版本/files/etc/rc.local
# Put your custom commands here that should be executed once
# the system init finished. By default this file does nothing.
init_gpio.sh
sleep 3
feed_dog.sh &
#gps_start.sh &
reconnect.sh &
#cgps_start.sh &
#探针启动
tysniffer -f /mnt/tysniffer.json -b true
exit 0
